# os/linux/configs/enable-ipv6-loopback.service
[Unit]
Description=Enable IPv6 loopback for the Adaptive stack

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -lc 'ip link set lo multicast on && ip route add ff01::0/16 dev lo'

[Install]
WantedBy=getty.target
